# src/deploy/pdf_validator.service
# /etc/systemd/system/pdf_validator.service
[Unit]
Description=PDF Validator (SigViewer web interface)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/pdf_validator

# --- Cache: persistent, writable, SHARED with the cron prefill ----------------
# CacheDirectory creates /var/cache/pdf_validator owned by the service user. It
# is writable even under ProtectSystem=strict (no ReadWritePaths needed) and,
# unlike RuntimeDirectory, it SURVIVES restarts. It is a real host directory
# (not namespaced), so the cron job writes the same files the service reads.
CacheDirectory=pdf_validator
CacheDirectoryMode=0750

# One fixed cache path for BOTH the service and cron. This explicit value wins
# over any $HOME/$XDG guessing, so sandboxing can't send them to different dirs.
Environment=SIGVIEWER_CACHE_DIR=/var/cache/pdf_validator/sigviewer
Environment=SIGVIEWER_LANG=en

# --- gunicorn's own runtime state ($HOME/.gunicorn) ---------------------------
# /var/www (www-data's home) is read-only under ProtectSystem=strict, so point
# HOME at a throwaway /run dir. RuntimeDirectory is correct HERE (state is
# disposable) — it must NOT hold the cache, which is why the cache now has its
# own SIGVIEWER_CACHE_DIR above instead of following HOME.
RuntimeDirectory=pdf_validator
Environment=HOME=/run/pdf_validator

ExecStart=/opt/pdf_validator/venv/bin/gunicorn \
    --workers 1 --threads 8 \
    --bind 127.0.0.1:8000 \
    --timeout 300 \
    webapp:app
Restart=on-failure
RestartSec=2

# --- Hardening (CacheDirectory + RuntimeDirectory stay writable) --------------
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
